# crontab
# Run the maintainer_db.py script every 10 minutes
*/2 * * * * su appuser && cd /app/request_handler && ls && /usr/local/bin/python maintainer_db.py >> /var/log/cron.log 2>&1
