# crontab
# Run the maintainer_db.py script every 10 minutes
*/10 * * * * appuser python /app/maintainer_db.py >> /var/log/cron.log 2>&1
